# knowledge-base/scripts/kb-sync.service
[Unit]
Description=Knowledge Base Auto Sync

[Service]
Type=oneshot
ExecStart=%h/andro-tools/knowledge-base/scripts/autocommit.sh
